# druid/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:druid]
command=/opt/druid/ptrun %(ENV_DRUID_FORM)s %(ENV_DRUID_OPERATION)s
